# app/supervisor.conf.d/novnc.conf
[program:novnc]
command="/noVNC/utils/novnc_proxy" --listen "0.0.0.0:%(ENV_NOVNC_PORT)s"
autorestart=true