# app/supervisor.conf.d/novnc.conf
[program:novnc]
command="/noVNC/utils/novnc_proxy" --listen "0.0.0.0:5000"
autorestart=true