# app/supervisor.conf.d/roscore.conf
[program:roscore]
command=roscore
autorestart=true