# app/supervisor.conf.d/roscore.conf
[program:roscore]
command=source /opt/ros/"%(ENV_ROS_DISTRO)s"/setup.bash && roscore
autorestart=true
user=ros