# app/supervisor.conf.d/roscore.conf
[program:novnc]
command=roscore
autorestart=true